# src/etc/bash/adas_laser_pub.service
[Unit]

Description=adas_laser_pub
After=network.target

[Service]
User=nvidia
ExecStart=/bin/bash /home/nvidia/mdl_ws/src/etc/bash/adas_laser_pub.bash
Restart=on-failure
RestartSec=5s

[Install]

WantedBy=multi-user.target
